# packages/debseed/debseed-0.1.3-py3-none-any.whl/debseed/assets/late-command-assets/target/etc/systemd/system/debseed-firstboot.service
[Unit]
Before=systemd-user-sessions.service
Wants=network-online.target
After=network-online.target
ConditionPathExists=/var/local/debseed-firstboot/required.top

[Service]
Type=oneshot
ExecStart=/var/local/debseed-firstboot/main

[Install]
WantedBy=multi-user.target
